# salt/elife-dashboard/config/etc-init-uwsgi-elife-dashboard-2.conf
description "dashboard_2 uWSGI server"
start on (local-filesystems and net-device-up IFACE!=lo)
respawn
exec /srv/elife-dashboard/dashboard_2/venv/bin/uwsgi --die-on-term --ini /srv/elife-dashboard/dashboard_2/uwsgi.ini
